# deployment/supervisord.conf
# 阿尔兹海默症诊断系统 - Supervisor进程管理配置
# Supervisor是一个客户端/服务器系统，用于监控和控制UNIX进程

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
childlogdir=/app/logs
user=root

[program:uvicorn]
command=uvicorn src.app:app --host 0.0.0.0 --port 8888 --workers 4
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/uvicorn.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopwaitsecs=60
startretries=3
exitcodes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/nginx.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopsignal=QUIT

[group:alzheimer]
programs=uvicorn,nginx
priority=999